# com.distiller.bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.distiller.bot</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/Henry/Project/Distiller/scripts/run_bot.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/Henry/Project/Distiller</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/Henry/Project/Distiller/venv/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/Henry</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/Henry/Project/Distiller/logs/bot.out</string>

    <key>StandardErrorPath</key>
    <string>/Users/Henry/Project/Distiller/logs/bot.err</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>
</dict>
</plist>
