# com.distiller.diffords.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.distiller.diffords</string>

    <key>ProgramArguments</key>
    <array>
        <string>/path/to/distiller/scripts/run_diffords.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>4</integer>
        <key>Minute</key>
        <integer>0</integer>
        <key>Weekday</key>
        <integer>0</integer>
    </dict>

    <key>WorkingDirectory</key>
    <string>/path/to/distiller</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/home/USERNAME/.local/bin:/path/to/distiller/.venv/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/home/USERNAME</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/path/to/distiller/logs/diffords.out</string>

    <key>StandardErrorPath</key>
    <string>/path/to/distiller/logs/diffords.err</string>

    <key>RunAtLoad</key>
    <false/>

    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
